# files/supervisord.conf
[unix_http_server]
file=/dev/shm/supervisor.sock

[supervisord]
logfile=/tmp/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=true
minfds=1024
minprocs=200
user=root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///dev/shm/supervisor.sock

[program:eventserver]
command=pio eventserver
directory=/root
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/pio-event.server.err.log
stdout_logfile=/var/log/pio-event.server.out.log
user=root

[program:recommender]
command=pio deploy
directory=/root/engines/recommender-incubating
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/pio.server.err.log
stdout_logfile=/var/log/pio.server.out.log
user=root
